# theme/crontab
fish=/home/philip/.local/bin/fish

20,40 5,16 * * * $fish -lc "theme"
*/20 6-7,17-20 * * * $fish -lc "theme"
0 8,21 * * * $fish -lc "theme"
